# config/chroot_local-includes/etc/skel/.config/autostart/fusesmb.desktop
[Desktop Entry]
Encoding=UTF-8
Version=1.0
Type=Application
Exec=/usr/bin/fusesmb ~/WindowsNetwork -o allow_other
Name=fusesmb
