# config/chroot_local-includes/etc/skel/.config/autostart/mountsmb.desktop
[Desktop Entry]
Encoding=UTF-8
Version=1.0
Type=Application
Exec=/usr/local/bin/mount-fusesmb.sh
Name=Mount WindowsNetwork
